# debian/nodedaemon-walrus.postinst
#!/bin/bash -e

if [ "$1" = "configure" ]; then
    chown -R luhya:luhya /usr/local/nodedaemon || true
    chmod -R 777 /usr/local/nodedaemon || true

    ln -s /usr/local/nodedaemon/walrus/walrus_status_publisher /usr/bin/walrus_status_publisher || true

    systemctl enable supervisor || true
    systemctl start supervisor  || true
fi
