# debian/nodedaemon-nc.postinst
#!/bin/bash -e

if [ "$1" = "configure" ]; then
    chown -R luhya:luhya /usr/local/nodedaemon || true
    chmod -R 777 /usr/local/nodedaemon || true

    ln -s /usr/local/nodedaemon/nc/nc_cmd_consumer     /usr/bin/nc_cmd_consumer || true
    ln -s /usr/local/nodedaemon/nc/nc_status_publisher /usr/bin/nc_status_publisher || true

    systemctl enable supervisor
    systemctl start supervisor
fi
